# user_scripts/update_dusky/update_checker/service/update_checker.timer
[Unit]
Description=Dusky Update Checker Timer (Every 20m)

[Timer]
# Wait 5 minutes after boot/login before the first check
# (Prevents slowing down the startup process)
OnBootSec=5min

# Run 20 minutes after the last run finished
OnUnitActiveSec=20min

# IMPORTANT for distributed scripts:
# Adds a random delay of 0-5 minutes to prevent all users
# from hitting GitHub at the exact same second.
RandomizedDelaySec=300

# If the computer was off during the last scheduled run,
# run immediately upon booting up.
Persistent=true

[Install]
WantedBy=timers.target
